# medicore_server/com.medicore.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.medicore.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/eye/medicore_server/medicore-server</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Applications/eye/medicore_server</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>DB_HOST</key>
        <string>localhost</string>
        <key>DB_PORT</key>
        <string>5432</string>
        <key>DB_USER</key>
        <string>medicore</string>
        <key>DB_PASSWORD</key>
        <string>medicore_secure_2025</string>
        <key>DB_NAME</key>
        <string>medicore_db</string>
        <key>DB_SSLMODE</key>
        <string>disable</string>
        <key>SERVER_PORT</key>
        <string>50052</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Applications/eye/medicore_server/logs/server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Applications/eye/medicore_server/logs/server_error.log</string>
</dict>
</plist>
